# systemd/ob-heartbeat.service
[Unit]
Description=Open Bastion Heartbeat
Documentation=man:ob-heartbeat(8)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/ob-heartbeat
User=root
Group=root

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

# ob-heartbeat refreshes the access token, rewriting it atomically (mktemp +
# rename). The token now lives under /var/lib/open-bastion (runtime state, FHS)
# — already writable below — so /etc stays fully read-only under
# ProtectSystem=strict (config is read-only at runtime, as it should be).
# openbastion.conf is read from the read-only /etc; nothing in /etc is written.
ReadWritePaths=/var/lib/open-bastion

# Network access required
PrivateNetwork=no
