# systemd/ob-heartbeat.service
[Unit]
Description=Open Bastion Heartbeat
Documentation=man:ob-heartbeat(8)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/ob-heartbeat
User=root
Group=root

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

# Allow reading config and token files
ReadOnlyPaths=/etc/open-bastion/openbastion.conf
ReadOnlyPaths=/etc/open-bastion/token
ReadWritePaths=/var/lib/open-bastion

# Network access required
PrivateNetwork=no
