# debian/postinst
#!/bin/sh
set -e

. /usr/share/debconf/confmodule

CONF_FILE="/etc/security/pam_llng.conf"
BACKUP_DIR="/var/backups/libpam-llng"
PAM_SSHD="/etc/pam.d/sshd"
PAM_SUDO="/etc/pam.d/sudo"

# Backup a PAM file if not already backed up
backup_pam_file() {
    file="$1"
    name="$2"
    if [ -f "$file" ] && [ ! -f "$BACKUP_DIR/${name}.orig" ]; then
        mkdir -p "$BACKUP_DIR"
        cp -a "$file" "$BACKUP_DIR/${name}.orig"
    fi
}

# Generate PAM configuration for sshd
generate_pam_sshd() {
    mode="$1"
    server_group="$2"
    cat << EOF
# /etc/pam.d/sshd - Configured by libpam-llng
# Mode: $mode
# Server group: $server_group
#
EOF
    case "$mode" in
        mode-a)
            cat << EOF
# AUTHENTICATION: Only LLNG tokens accepted
auth       sufficient   pam_llng.so server_group=$server_group
auth       required     pam_deny.so
EOF
            ;;
        mode-b|mode-d)
            cat << EOF
# AUTHENTICATION: LLNG token OR Unix password
auth       sufficient   pam_llng.so server_group=$server_group
auth       sufficient   pam_unix.so nullok try_first_pass
auth       required     pam_deny.so
EOF
            ;;
        mode-c)
            cat << 'EOF'
# AUTHENTICATION: SSH keys only (handled by sshd, not PAM)
auth       required     pam_permit.so
EOF
            ;;
    esac
    cat << EOF

# AUTHORIZATION: LLNG checks if user can access this server
account    required     pam_llng.so server_group=$server_group
account    required     pam_unix.so

# SESSION
session    required     pam_unix.so
EOF
}

# Generate PAM configuration for sudo
generate_pam_sudo() {
    mode="$1"
    server_group="$2"
    cat << EOF
# /etc/pam.d/sudo - Configured by libpam-llng
# Mode: $mode
# Server group: $server_group
#
EOF
    case "$mode" in
        mode-a)
            cat << EOF
# AUTHENTICATION: Only LLNG tokens accepted
auth       sufficient   pam_llng.so server_group=$server_group
auth       required     pam_deny.so
EOF
            ;;
        mode-b|mode-d)
            cat << EOF
# AUTHENTICATION: LLNG token OR Unix password
auth       sufficient   pam_llng.so server_group=$server_group
auth       sufficient   pam_unix.so nullok try_first_pass
auth       required     pam_deny.so
EOF
            ;;
        mode-c)
            cat << 'EOF'
# AUTHENTICATION: Permit (for SSH key scenarios)
auth       required     pam_permit.so
EOF
            ;;
    esac
    cat << EOF

# AUTHORIZATION
account    required     pam_llng.so server_group=$server_group
account    required     pam_unix.so

# SESSION
session    required     pam_unix.so
EOF
}

# Generate pam_llng.conf
generate_llng_conf() {
    portal_url="$1"
    client_id="$2"
    client_secret="$3"
    server_group="$4"

    cat << EOF
# /etc/security/pam_llng.conf
# Configured by libpam-llng debconf
# Run 'dpkg-reconfigure libpam-llng' to modify

# LemonLDAP::NG portal URL
portal_url = $portal_url

# OIDC client credentials
client_id = $client_id
client_secret = $client_secret

# Server authentication
server_token_file = /etc/security/pam_llng.token
server_group = $server_group

# HTTP settings
timeout = 10
verify_ssl = true

# Cache settings
cache_enabled = true
cache_dir = /var/cache/pam_llng
cache_ttl = 300

# Logging: error, warn, info, debug
log_level = warn
EOF
}

case "$1" in
    configure)
        db_get libpam-llng/pam-mode
        MODE="$RET"

        if [ "$MODE" != "none" ] && [ -n "$MODE" ]; then
            # Get debconf values
            db_get libpam-llng/portal-url
            PORTAL_URL="$RET"
            db_get libpam-llng/client-id
            CLIENT_ID="$RET"
            db_get libpam-llng/client-secret
            CLIENT_SECRET="$RET"
            db_get libpam-llng/server-group
            SERVER_GROUP="$RET"
            db_get libpam-llng/configure-sudo
            CONFIGURE_SUDO="$RET"
            db_get libpam-llng/sudo-server-group
            SUDO_SERVER_GROUP="$RET"

            # Use main server group for sudo if not specified
            if [ -z "$SUDO_SERVER_GROUP" ]; then
                SUDO_SERVER_GROUP="$SERVER_GROUP"
            fi

            # Generate pam_llng.conf if portal URL is provided
            if [ -n "$PORTAL_URL" ]; then
                generate_llng_conf "$PORTAL_URL" "$CLIENT_ID" "$CLIENT_SECRET" "$SERVER_GROUP" > "$CONF_FILE"
                chmod 600 "$CONF_FILE"
                chown root:root "$CONF_FILE"
            fi

            # Backup and configure sshd
            backup_pam_file "$PAM_SSHD" "sshd"
            generate_pam_sshd "$MODE" "$SERVER_GROUP" > "$PAM_SSHD"

            # Configure sudo if requested
            if [ "$CONFIGURE_SUDO" = "true" ]; then
                backup_pam_file "$PAM_SUDO" "sudo"
                generate_pam_sudo "$MODE" "$SUDO_SERVER_GROUP" > "$PAM_SUDO"
            fi

            # Create cache directory
            mkdir -p /var/cache/pam_llng
            chmod 700 /var/cache/pam_llng

            # Create stats directory
            mkdir -p /var/lib/pam_llng
            chmod 700 /var/lib/pam_llng

            # Enable and start heartbeat timer
            if [ -d /run/systemd/system ]; then
                systemctl daemon-reload || true
                systemctl enable pam-llng-heartbeat.timer || true
                systemctl start pam-llng-heartbeat.timer || true
            fi

            # Clear secret from debconf database for security
            db_set libpam-llng/client-secret ""
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
